# resources/chv-cmdserver.service
[Unit]
Description=CHV Command Server
After=chv-guestinit.service

[Service]
Type=simple
ExecStart=/usr/local/bin/chv-cmdserver
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
